# Twitter_IconChange/systemd_sample/changeicon.timer
[Unit]
Description= changeicon timer

[Timer]
OnUnitActiveSec=360m

[Install]
WantedBy=timers.target