# MyIconChanger/systemd_sample/changeicon.service
[Unit]
Description= ChangeIcon Service

[Service]
Type= oneshot
ExecStart= /usr/bin/php /MyIconChanger/change.php